# meta-bsp/recipes-myir/fac-burn-emmc-core/fac-burn-emmc-core/fac-burn-emmc.service
[Unit]
Description=auto run script after psplash-quit service
After=psplash-quit.service

[Service]
Type=oneshot
ExecStart=/home/root/burn_emmc.sh


[Install]
WantedBy=multi-user.target
