# meta-sdk/recipes-myir/qt-demo/qt-demo/myir.service
[Unit]
Description=myir hmi v2.0
After=getty@tty1.service

[Service]
Type=forking
ExecStart=/bin/sh /usr/bin/start.sh

[Install]
WantedBy=multi-user.target
Alias=myir.service

